# worker/latch-agent-bridge.service
[Unit]
Description=Latch approval-planning bridge for OpenClaw worker
After=network-online.target tailscaled.service docker.service
Wants=network-online.target

[Service]
# Type=notify + WatchdogSec: the bridge sends WATCHDOG=1 each tick (see sd_notify in the
# script). If the tick loop stalls (a hang no timeout catches), systemd force-restarts the
# service instead of leaving it silently wedged and not consuming any work.
Type=notify
WatchdogSec=180
EnvironmentFile=/etc/latch-agent-bridge.env
ExecStart=/usr/local/bin/latch-agent-bridge.py --state-path /var/lib/latch-agent-bridge/state.json --source-cache-path /var/lib/latch-agent-bridge/source-notes.json
Restart=always
RestartSec=10
DynamicUser=true
StateDirectory=latch-agent-bridge
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/latch-agent-bridge

[Install]
WantedBy=multi-user.target
